# root-overlay/etc/cron.hourly/ntpdate
#!/bin/bash

ntpdate -u pool.ntp.org -s

